# .assets/scripts/stsd.service
[Unit]
Description=Site Transfer Software Daemon
After=syslog.target network.target
 
[Service]
Type=simple
User=dsmgr
Group=data
RestartSec=5
ExecStart=/apps/transfer/bin/sts -loop -conf /apps/transfer/conf/sts.yaml
WatchdogSec=30s
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
